# artifacts/systemd/backup-system.service
[Unit]
Description=System backup via Borg

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-system.sh

