# artifacts/systemd/backup-system.timer
[Unit]
Description=Daily system backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

